# docker/docker_jit_monitor/etc/notify-email@.service
[Unit]
Description=Send email

[Service]
Type=oneshot
ExecStart=/usr/bin/bash -c '/home/servo_ci/.local/bin/servo_ci_runner_send_mail.sh "%i"'

